# gazebo/gazebo/supervisord.conf
[supervisord]
nodaemon=true

[program:gzserver]
command = /bin/bash -c "gazebo worlds/pr2.world"
