# docker/gazebo/web/supervisord.conf
[supervisord]
nodaemon=true

[program:gzserver]
command = /bin/bash -c "gzserver worlds/pr2.world"

[program:gzclient]
command = /bin/bash -c "./~/gzweb/start_gzweb.sh"
